# openkds.service
# Ce fichier est une référence — install.sh l'écrit directement dans /etc/systemd/system/
[Unit]
Description=OpenKDS Restaurant System
After=network.target openkds-hotspot.service
Wants=openkds-hotspot.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=openkds
Environment=OPENKDS_DATA_DIR=/opt/openkds/data
ExecStart=/opt/openkds/venv/bin/openkds
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
